# exam-registration-service/supervisor/response-saga-item-consumer.conf
[program:response-saga-item-consumer]
command=php /var/www/project/bin/console messenger:consume response_saga_item_transport --time-limit=3600
user=root
startsecs=0
autostart=true
autorestart=true
startretries=10