# email-service/supervisor/student-updated-consumer.conf
[program:student-updated-consumer]
command=php /var/www/project/bin/console messenger:consume student_update --time-limit=3600
user=root
startsecs=0
autostart=true
autorestart=true
startretries=10